# SWOOLE/supervisord.conf
[supervisord]
nodaemon=true
[program:swooleDistributed]
command=php /wwwroot/bin/start_swoole_server.php start -d
autostart=true
autorestart=true
redirect_stderr=true